# jtek-bot.timer
[Unit]
Description=JTEK/VOO/SGOV Bot — daily at 4:00 PM ET (20:00 UTC, after US close)

[Timer]
OnCalendar=*-*-* 20:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
